# DOP_automation_2019/roles/result/files/result.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/usr/bin/node /root/result/server.js

[Install]
WantedBy=multi-user.target